# static/files/distr/raspberry-pi/home/.config/autostart/autostart.desktop
[Desktop Entry]
Type=Application
Name=Autostart
Exec=/bin/bash /home/pi/.config/autostart/autostart.sh

