# static/files/distr/raspberry-pi/autostart/autostart.desktop
[Desktop Entry]
Type=Application
Name=SystemInfo
Exec=/bin/bash /home/pi/.config/autostart/autostart.sh

